# packages-travelping/mand/files/mand.init
#!/bin/sh /etc/rc.common

START=80
STOP=10

USE_PROCD=1
PROG=/usr/sbin/mand

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}
