# packages-travelping/mand-cfgd/files/mand-cfgd.init
#!/bin/sh /etc/rc.common

START=90
STOP=10

USE_PROCD=1
PROG=/usr/bin/mand-cfgd

start_service() {
	procd_open_instance
	procd_set_param command "$PROG" -f
	procd_close_instance
}
